# haos-overlay/buildroot-external/rootfs-overlay/etc/NetworkManager/dispatcher.d/99-hassos-dns-cn
#!/bin/sh
# NetworkManager dispatcher: end0 连接 up 时立即配置 host DNS，再后台配置 supervisor DNS
# 作为 hassos-dns-cn-init.service 的双保险：service 失败后，dispatcher 仍会在
# 网络恢复时尝试重新配置。
#
# Build #49 关键修复: 不再硬编码 192.168.1.1
# 不同家庭 IP 网段不同 (192.168.0.1/192.168.1.1/10.0.0.1/192.168.2.1 等)
# 用运行时检测的默认网关 IP, 适配任意家庭网段

INTERFACE="$1"
ACTION="$2"

[ "$INTERFACE" = "end0" ] || exit 0
[ "$ACTION" = "up" ] || [ "$ACTION" = "dhcp4-change" ] || exit 0

# 关键 (Build #49): 运行时自动检测默认网关 IP
# 优先级: ip route > nmcli IP4.GATEWAY > nmcli IP4.DNS > 兜底 192.168.1.1
DNS_SERVERS=""
if command -v ip >/dev/null 2>&1; then
    DNS_SERVERS=$(ip -4 route show default 2>/dev/null | awk '/^default/ {for(i=1;i<=NF;i++) if ($i=="via") {print $(i+1); exit}}' | head -1)
fi
if [ -z "$DNS_SERVERS" ] && command -v nmcli >/dev/null 2>&1; then
    DNS_SERVERS=$(nmcli -t -f IP4.GATEWAY device show end0 2>/dev/null | head -1)
fi
if [ -z "$DNS_SERVERS" ] && command -v nmcli >/dev/null 2>&1; then
    DNS_SERVERS=$(nmcli -t -f IP4.DNS device show end0 2>/dev/null | head -1 | cut -d: -f1)
fi
if [ -z "$DNS_SERVERS" ]; then
    DNS_SERVERS="192.168.1.1"
    logger -t hassos-dns-cn-dispatcher "WARNING: 无法自动检测默认网关, 使用兜底 192.168.1.1"
fi
logger -t hassos-dns-cn-dispatcher "检测到默认网关: $DNS_SERVERS (host DNS)"

# 立即配置 host DNS
# Build #50: 指向 hassio_dns (172.30.32.3), 由 hassio_dns 的 Corefile 多上游故障转移
# 处理实际解析。不再指向路由器 (实测 192.168.1.1:53 经常 i/o timeout)。
if command -v resolvectl >/dev/null 2>&1; then
    resolvectl dns end0 172.30.32.3 2>/dev/null || true
fi

if command -v nmcli >/dev/null 2>&1; then
    CONN_NAME=$(nmcli -t -f NAME,DEVICE connection show --active 2>/dev/null | grep ':end0$' | cut -d: -f1 | head -1)
    if [ -n "$CONN_NAME" ]; then
        nmcli connection modify "$CONN_NAME" ipv4.dns "172.30.32.3" ipv4.ignore-auto-dns yes 2>/dev/null || true
        nmcli connection modify "$CONN_NAME" ipv6.ignore-auto-dns yes 2>/dev/null || true
    fi
fi

# 后台执行 supervisor DNS 配置
/usr/libexec/hassos-dns-cn-init end0 &

exit 0
