# haos-overlay/buildroot-external/rootfs-overlay/etc/NetworkManager/dispatcher.d/99-hassos-dns-cn
#!/bin/sh
# NetworkManager dispatcher: end0 连接 up 时立即把 host DNS 指向 hassio_dns,
# 再后台调用 hassos-dns-cn-init 配置 supervisor DNS (Cloudflare DoT -> 公共 DNS)。
#
# Build #51 设计 (关键 loop 修复):
#   - host DNS 指向 hassio_dns (172.30.32.3); 实际解析由 CoreDNS 的 '.:53' ->
#     dns://127.0.0.11 (Docker 内嵌) -> daemon.json 的 dns (公共 DNS) 完成, 绝不直连网关。
#   - 公共 DNS 是全局固定 IP, 与家庭网段/网关无关, 天然适配任意路由器。
#   - 绝不在 CoreDNS forward 中直连路由器网关 IP (会触发 plugin/loop FATAL)。

INTERFACE="$1"
ACTION="$2"

[ "$INTERFACE" = "end0" ] || exit 0
[ "$ACTION" = "up" ] || [ "$ACTION" = "dhcp4-change" ] || exit 0

logger -t hassos-dns-cn-dispatcher "end0 $ACTION: 配置 host DNS -> hassio_dns (172.30.32.3)"

# 立即配置 host DNS (best-effort, 主脚本 hassos-dns-cn-init 会带重试再做一次)
if command -v resolvectl >/dev/null 2>&1; then
    resolvectl dns end0 172.30.32.3 2>/dev/null || true
fi

if command -v nmcli >/dev/null 2>&1; then
    CONN_NAME=$(nmcli -t -f NAME,DEVICE connection show --active 2>/dev/null | grep ':end0$' | cut -d: -f1 | head -1)
    if [ -n "$CONN_NAME" ]; then
        nmcli connection modify "$CONN_NAME" ipv4.dns "172.30.32.3" ipv4.ignore-auto-dns yes 2>/dev/null || true
        nmcli connection modify "$CONN_NAME" ipv6.ignore-auto-dns yes 2>/dev/null || true
    fi
fi

# 后台执行 supervisor DNS 配置 (带 flock 防并发)
/usr/libexec/hassos-dns-cn-init end0 &

exit 0
