# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/haos-log-capture.service
[Unit]
Description=HAOS log capture (dmesg/journal/pstore to /mnt/data/sos-logs)
After=local-fs.target systemd-journald.service
ConditionPathExists=/mnt/data

[Service]
Type=oneshot
ExecStart=/usr/libexec/haos-log-capture
RemainAfterExit=no
