# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/haos-log-capture.timer
[Unit]
Description=Periodic HAOS log capture

[Timer]
OnBootSec=15
OnUnitActiveSec=60
Unit=haos-log-capture.service

[Install]
WantedBy=timers.target
