# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/haos-ensure-files.service
[Unit]
Description=Ensure critical /etc files exist for HAOS
DefaultDependencies=no
Before=hassos-overlay.service hassos-bind.target
# Do NOT use After=local-fs.target — it creates an ordering cycle:
# local-fs.target → bind-mounts → hassos-overlay.service → haos-ensure-files.service → local-fs.target
# Instead, rely on the overlay partition mount (mnt-overlay.mount) being ready
After=mnt-overlay.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/haos-ensure-files

[Install]
WantedBy=hassos-bind.target
