# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-dns-cn-init-supervisor-recheck.timer
[Unit]
Description=HAOS-CN: Periodic supervisor recheck (every 5 min)
# 配合 hassos-dns-cn-init-supervisor-recheck.service 使用
# 直到 supervisor 真正检测到网络为止

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
# 启动后 5 分钟开始，每 5 分钟执行
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
