# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-dns-cn-init.service
[Unit]
Description=HAOS-CN: Configure hassio_dns with Chinese DNS servers (disable CloudFlare DoT fallback)
Documentation=https://developers.home-assistant.io/docs/api/supervisor/endpoints
# 必须等网络就绪 + supervisor 启动完成。Build #34 缺少 network-online.target，
# 导致 64s 本服务启动时 supervisor.sock 刚创建但 IPv4 尚未获取，curl 失败。
After=network-online.target hassos-supervisor.service
Wants=network-online.target hassos-supervisor.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/hassos-dns-cn-init
# 脚本内部最多等 240s (120s supervisor + 120s 网络)，留出余量
TimeoutStartSec=300s
# 失败时重试：DHCP 可能多次超时（如交换机 STP），supervisor 也可能晚于本服务启动
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
