# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-supervisor.service.d/resilience.conf
[Unit]
# RK3399 UBB: Give supervisor more retries before giving up
# StartLimitIntervalSec and StartLimitBurst belong in [Unit] section
StartLimitIntervalSec=10m
StartLimitBurst=10

[Service]
RestartSec=10s
TimeoutStartSec=300

# Standard output/error for debugging
StandardOutput=journal+console
StandardError=journal+console

# Environment file for supervisor (fixed: was ENVIRONMENT_FILE)
EnvironmentFile=-/run/dockerd.env
