# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-supervisor.service.d/resilience.conf
[Service]
# RK3399 UBB: Give supervisor more time and retries
# The board has limited RAM (2-4GB) and slower eMMC
StartLimitIntervalSec=10m
StartLimitBurst=10
RestartSec=10s
TimeoutStartSec=300

# Standard output/error for debugging
StandardOutput=journal+console
StandardError=journal+console

# Environment for supervisor
ENVIRONMENT_FILE=/run/dockerd.env
