# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/docker.service.d/failure-override.conf
[Unit]
# Override upstream failure.conf to prevent system reboot on docker failure
# The upstream has FailureAction=reboot which causes 60s reboot loops
# when docker fails to start. We disable this for debugging.
FailureAction=none
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
# Don't reboot the system if docker fails
# The supervisor will retry automatically
TimeoutStartSec=300
