# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/docker.service.d/override.conf
[Service]
# RK3399 cgroup v2 fix - use systemd cgroup driver
Delegate=yes
CPUAccounting=yes
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes
# Ensure docker can start even if some mounts are delayed
TimeoutStartSec=300
TimeoutStopSec=30
# Prevent OOM kill for docker
OOMScoreAdjust=-500
