# haos-overlay/buildroot-external/rootfs-overlay/usr/lib/systemd/system/docker.service.d/override.conf
[Service]
Delegate=yes
CPUAccounting=yes
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes
